# sv/tasksy/run
#! /usr/bin/env sh
DISPLAY=:0
snooze -H'*' -M/2 task sync
