# sv/syncthing/run
#! /usr/bin/env fish

syncthing serve
#--config=/config/syncthing --data=/config/syncthing/db --audit=true --auditfile=/config/syncthing/audit.txt
