# sv/background/run
#! /usr/bin/env sh
DISPLAY=:0
feh --bg-fill --randomize ~/.backgrounds
