# sv/selfhost/run
#! /usr/bin/env fish

sleep 10s
/code/selfhosting/main