# sv/jackett/finish
#! /usr/bin/env sh

podman stop jackett
podman rm jackett