# sv/emacs_server/run
#! /usr/bin/env fish

emacs --with-x-toolkit=lucid --daemon
